# supervisord.conf
###############################################################################
#  supervisor — jalanin beberapa proses sekaligus dalam 1 container           #
#  File ini di-COPY ke /etc/supervisor/conf.d/hub.conf                        #
#  Program cloudflared TIDAK di sini — ditambah dinamis oleh start.sh         #
#  (cuma kalau CF_TUNNEL_TOKEN diisi), biar tes lokal gak error berisik.      #
###############################################################################

[program:sshd]
command=/usr/sbin/sshd -D -e
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:frps]
command=/usr/local/bin/frps -c /etc/frp/frps.toml
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=15

[program:sshwifty]
command=/usr/local/bin/sshwifty
environment=SSHWIFTY_CONFIG="/etc/sshwifty/sshwifty.conf.json"
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
